# isosource/chroot/var/lib/dpkg/info/apt.list
/var
/var/log
/var/log/apt
/var/cache
/var/cache/apt
/var/cache/apt/archives
/var/lib
/var/lib/apt
/var/lib/apt/lists
/etc
/etc/cron.daily
/etc/cron.daily/apt
/etc/logrotate.d
/etc/logrotate.d/apt
/etc/kernel
/etc/kernel/postinst.d
/etc/kernel/postinst.d/apt-auto-removal
/etc/apt
/etc/apt/apt.conf.d
/etc/apt/apt.conf.d/20changelog
/etc/apt/apt.conf.d/01-vendor-ubuntu
/etc/apt/apt.conf.d/01autoremove
