# isosource/chroot/var/lib/dpkg/info/plymouth.list
/etc
/etc/init.d
/etc/init.d/plymouth
/etc/init.d/plymouth-log
