# isosource/chroot/var/lib/dpkg/info/ufw.list
/etc
/etc/default
/etc/default/ufw
/etc/init.d
/etc/init.d/ufw
/etc/logrotate.d
/etc/logrotate.d/ufw
/etc/rsyslog.d
/etc/rsyslog.d/20-ufw.conf
/etc/ufw
/etc/ufw/sysctl.conf
