# isosource/chroot/etc/rc.local
#!/bin/sh -e

#load web10g module
/sbin/modprobe tcp_estats_nl

# set the web10g variables
/sbin/sysctl -w net.ipv4.tcp_estats=127
/sbin/sysctl -w net.ipv4.estats_delay=1500

#force the interfaces to get an ip address
/sbin/dhclient

# get/set  the current time
/usr/sbin/ntpd

exit 0
